# instances/default/projects/712f694b-097f-4ca4-b161-f299187ec27d/9fde99f6-5309-4ba7-98d8-bfcdd943e3b5/_default/nl.webbio.paperclip.ai-news-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>nl.webbio.paperclip.ai-news-cron</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/joerismits/.paperclip/instances/default/projects/712f694b-097f-4ca4-b161-f299187ec27d/9fde99f6-5309-4ba7-98d8-bfcdd943e3b5/_default/run-ai-news-cron.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/paperclip-ai-news.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/paperclip-ai-news.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
